# websites/atviriduomenys.lt/templates/spinta.socket.j2
[Unit]
Description=spinta socket

[Socket]
ListenStream=/run/spinta/socket
ListenStream=127.0.0.1:8000

[Install]
WantedBy=sockets.target
